# remote/systemd/uhrk-backend.service
[Unit]
Description=UHRK SX1303 telemetry backend
After=network-online.target
Wants=network-online.target
Conflicts=chirpstack-gateway-bridge.service

[Service]
Type=simple
User=uhrkgc
WorkingDirectory=/opt/uhrk/uhrk_site
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/uhrk/gc.env
ExecStart=/usr/bin/python3 /opt/uhrk/uhrk_site/uhrk_backend.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
