# remote/systemd/uhrk-node.service
[Unit]
Description=UHRK telemetry node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=uhrk
WorkingDirectory=/opt/uhrk
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/uhrk
EnvironmentFile=-/etc/uhrk/node.env
ExecStart=/opt/uhrk/node-venv/bin/python -m zenith_node.main
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
